# root/etc/systemd/system/tomcat8@.service.d/webtop.conf
[Unit]
Wants=postgresql.service
After=postgresql.service

[Service]
StandardOutput=null
StandardError=null
Restart=on-failure
