# root/etc/systemd/system/tomcat8@.service.d/webtop.conf
[Service]
StandardOutput=null
StandardError=null
Restart=on-failure
